# scripts/systemd/ebs-netwatch-publish.timer
[Unit]
Description=Run EBS Netwatch publish every 30 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min
Persistent=true
Unit=ebs-netwatch-publish.service

[Install]
WantedBy=timers.target
